# systemd_data/data_script.timer
[Unit]
Description=data script timer

[Timer]
OnBootSec=0min
OnCalendar=*:*:0/30
Unit=data_script.service

[Install]
WantedBy=multi-user.target
