# systemd_data/data_script.service
[Unit]
Description=Run the data script

[Service]
ExecStart=/Users/.../date/dateToSystemd.sh
